# deploy/com.user.xdl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS LaunchAgent: keeps the local engine running whenever you're logged in,
  bound to localhost, writing to ~/Downloads. Optional — skip it and use
  `xdl get <url>` for synchronous one-shots instead.

  Install:
    1. Replace __PROJECT_DIR__ and __HOME__ below (or use the sed line in
       deploy/macos-client.md, which does it for you).
    2. cp this file to ~/Library/LaunchAgents/com.user.xdl.plist
    3. launchctl load ~/Library/LaunchAgents/com.user.xdl.plist
    4. launchctl list | grep xdl   # confirm it's running
  Uninstall: launchctl unload ~/Library/LaunchAgents/com.user.xdl.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.xdl</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/.venv/bin/xdl</string>
        <string>serve</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>XDL_HOST</key>            <string>127.0.0.1</string>
        <key>XDL_PORT</key>            <string>8080</string>
        <key>XDL_DB</key>             <string>__HOME__/.local/share/xdl/ledger.db</string>
        <key>XDL_ARCHIVE</key>        <string>__HOME__/.local/share/xdl/archive.txt</string>
        <key>XDL_HUEY_DB</key>        <string>__HOME__/.local/share/xdl/huey.db</string>
        <key>XDL_DOWNLOAD_DIR</key>   <string>__HOME__/Downloads</string>
        <key>XDL_STAGING_DIR</key>    <string>__HOME__/.local/share/xdl/staging</string>
        <key>XDL_WORKERS</key>        <string>4</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/share/xdl/engine.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.local/share/xdl/engine.log</string>
</dict>
</plist>
